# deploy/runner-dashboard.service
# Single source of truth for the runner-dashboard systemd unit.
# deploy/setup.sh substitutes the @VAR@ placeholders below at install time.
# Substitution variables: @USER@, @HOME@, @DEPLOY_DIR@, @PYTHON_BIN@,
#   @MACHINE_NAME@, @NUM_RUNNERS@, @PORT@, @DISPLAY_NAME@,
#   @RUNNER_ALIASES@, @MACHINE_ROLE@, @SCHEDULE_CONFIG@.
# Hardening directives are static and MUST NOT be edited per machine.
# Tests in tests/test_systemd_template_parity.py enforce that this template
# is the only place hardening directives live.
[Unit]
Description=D-sorganization Runner Dashboard (@MACHINE_NAME@)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@DEPLOY_DIR@
ExecStartPre=@DEPLOY_DIR@/refresh-token.sh
ExecStart=@PYTHON_BIN@ @DEPLOY_DIR@/backend/server.py
Restart=always
RestartSec=5
Environment=GITHUB_ORG=D-sorganization
Environment=NUM_RUNNERS=@NUM_RUNNERS@
Environment=DASHBOARD_PORT=@PORT@
Environment=DISPLAY_NAME=@DISPLAY_NAME@
Environment=RUNNER_ALIASES=@RUNNER_ALIASES@
Environment=RUNNER_SCHEDULE_CONFIG=@SCHEDULE_CONFIG@
Environment=RUNNER_SCHEDULER_BIN=/usr/local/bin/runner-scheduler
Environment=MACHINE_ROLE=@MACHINE_ROLE@
Environment=HOME=@HOME@
Environment=PATH=/usr/lib/wsl/lib:@HOME@/.local/bin:@HOME@/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Secrets (GH_TOKEN, FLEET_NODES) are loaded from a file readable only by
# this user — not stored in this world-readable unit file.
# File: @HOME@/.config/runner-dashboard/env
# Contents:
#   GH_TOKEN=ghp_<classic PAT with scopes: repo, admin:org>
#   FLEET_NODES=Linux-Brick:http://100.x.x.a:8321,...   (once nodes are live)
EnvironmentFile=-@HOME@/.config/runner-dashboard/env

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RemoveIPC=true
# Allow the dashboard to read/write runner secrets and config from HOME.
ReadWritePaths=@HOME@/.config/runner-dashboard

[Install]
WantedBy=multi-user.target
